# installer/opensuse/sysmanage-agent.service
[Unit]
Description=SysManage Agent - System Management and Monitoring
Documentation=https://github.com/bceverly/sysmanage-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sysmanage-agent
Group=sysmanage-agent
WorkingDirectory=/opt/sysmanage-agent
ExecStart=/opt/sysmanage-agent/.venv/bin/python3 /opt/sysmanage-agent/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sysmanage-agent

# Security hardening — note that NoNewPrivileges is intentionally NOT
# set here.  The agent is designed around a "non-root user + sudoers
# NOPASSWD allowlist" privilege model (see sysmanage-agent.sudoers):
# the agent runs as ``sysmanage-agent``, but every privileged
# operation (zypper/rpm, systemctl, install, firewall-cmd, etc.) is
# invoked via ``sudo -n``, which the sudoers file restricts to a
# specific command allowlist.  NoNewPrivileges=true would make the
# kernel reject those sudo elevations regardless of the sudoers
# config, breaking every privileged operation the agent is designed
# to perform (mirror apply, package install, service restart,
# firewall config, etc.).  Hardening here is provided by the
# sudoers allowlist scope rather than the systemd
# ``NoNewPrivileges`` flag.
PrivateTmp=true
ProtectHome=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Environment
Environment="PYTHONUNBUFFERED=1"
Environment="SYSMANAGE_CONFIG=/etc/sysmanage-agent.yaml"
Environment="PATH=/opt/sysmanage-agent/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[Install]
WantedBy=multi-user.target
